# templates/ipa_cert_expiry_exporter.service.j2
[Unit]
Description=Update IPA certificate expiration metrics

[Service]
Type=oneshot
ExecStart={{ script_path }}
User=root